# templates/default/session.conf.erb
#!upstart
description "iCornerstone Session Service"
author      "John Manero <jmanero@dyn.com>"

start on startup
stop on shutdown

setuid <%= node["session_service"]["user"] %>
console log
chdir <%= node["session_service"]["home"] %>
exec <%= node["session_service"]["home"] %>/server.js
